# ops/systemd/navalivay-server.service
[Unit]
Description=NAVALIVAY Mini App Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/NAVALIVAY/server
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10
EnvironmentFile=/var/www/NAVALIVAY/server/.env

# Безопасность (раскомментируйте после создания пользователя)
# User=www-data
# Group=www-data

# Лимиты
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=navalivay-server

# Защита от зависаний
StartLimitInterval=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
